# root/etc/services.d/octoprint/run
#!/usr/bin/execlineb -P
s6-setuidgid octoprint
/usr/local/bin/octoprint serve --host 0.0.0.0
